# 6.DeployApp/app.service
[Unit]
Description=My Flask Application
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/usr/bin/python3 /home/ubuntu/app.py
Restart=always
Environment=FLASK_ENV=production

[Install]
WantedBy=multi-user.target
